# systemd/tplink-router-reboot.service
[Unit]
Description=Reboot the router

[Service]
ExecStart=tplink-router-reboot
